# docker/assets/rc.local
#!/bin/sh
echo "Starting rc.local"
BBB_HOST=`cat /etc/bigbluebutton/bbb-web.properties | grep 'bigbluebutton.web.serverURL='  | awk -F '://' '{print $2}'`
THIS_HOST=`hostname -f`

umount /etc/letsencrypt/live/`hostname -f` || true
rm -rf /etc/letsencrypt
mkdir /etc/letsencrypt/live/`hostname -f` -p
mount --bind /local/certs/ /etc/letsencrypt/live/`hostname -f`


# now haproxy serve port 443, and it expects the cert in /etc/haproxy/certbundle.pem
#copied from https://github.com/bigbluebutton/bbb-install/blob/16825846cb230ead230f00c9d3064b8d6dcee0bf/bbb-install-2.7.sh#L761
cat /local/certs/fullchain.pem /local/certs/privkey.pem > /etc/haproxy/certbundle.pem.new
chown root:haproxy /etc/haproxy/certbundle.pem.new
chmod 0640 /etc/haproxy/certbundle.pem.new
mv /etc/haproxy/certbundle.pem.new /etc/haproxy/certbundle.pem


if [ "$BBB_HOST" != "$THIS_HOST" ] ; then
    sed -i 's/'$BBB_HOST'/'$THIS_HOST'/g' /etc/nginx/sites-available/bigbluebutton
    sed -i 's/'$BBB_HOST'/'$THIS_HOST'/g' /usr/share/bigbluebutton/nginx/sip.nginx
    /usr/bin/bbb-conf --setip "$THIS_HOST"
fi;

LOCAL_IP=`hostname -I | awk '{print $1}'`
sed -i 's/data="local_ip_v4=.*"/data="local_ip_v4='$LOCAL_IP'"/g' /opt/freeswitch/conf/vars.xml

sed -i 's/bbb-docker-build.bbb.imdt.dev/'`hostname -f`'/g' /etc/nginx/sites-available/bigbluebutton
# Tweak COTURN configs with correct hostname and IP
sudo sed -i "s/172.17.0.2/$LOCAL_IP/g" /etc/turnserver.conf
sudo sed -i 's/bbb-docker-build.bbb.imdt.dev/'`hostname -f`'/g' /etc/turnserver.conf
sudo sed -i 's/bbb-docker-build.bbb.imdt.dev/'`hostname -f`'/g' /etc/bigbluebutton/turn-stun-servers.xml
sudo sed -i "s/172.17.0.2/$LOCAL_IP/g" /etc/haproxy/haproxy.cfg
systemctl reload haproxy

# Set configs with correct IP
sudo sed -i "s/172.17.0.2/$LOCAL_IP/g" /etc/bigbluebutton/bbb-webrtc-sfu/production.yml
sudo sed -i "s/172.17.0.2/$LOCAL_IP/g" /usr/share/bbb-apps-akka/conf/application.conf
sudo sed -i "s/172.17.0.2/$LOCAL_IP/g" /usr/share/bbb-web/WEB-INF/classes/bigbluebutton.properties

if [ -f /local/certs/bbb-dev-ca.crt ]; then
    echo "Trusting root CA";
    mkdir /usr/local/share/ca-certificates/bbb-dev/
    cp /local/certs/bbb-dev-ca.crt /usr/local/share/ca-certificates/bbb-dev/
    update-ca-certificates
fi;

# https://docs.bigbluebutton.org/2.4/dev.html#switch-nginx-to-redirect-requests-to-meteor
cat /usr/share/bigbluebutton/nginx/bbb-html5.nginx  | grep "for production" -v > /tmp/bbb-html5.nginx ; cat /tmp/bbb-html5.nginx > /usr/share/bigbluebutton/nginx/bbb-html5.nginx ; rm /tmp/bbb-html5.nginx
cat /usr/share/bigbluebutton/nginx/bbb-html5.nginx | sed 's|# proxy_pass http://127.0.0.1:4100; # use for development|proxy_pass http://127.0.0.1:4100; # use for development|g' > /tmp/bbb-html5.nginx ; cat /tmp/bbb-html5.nginx > /usr/share/bigbluebutton/nginx/bbb-html5.nginx ; rm /tmp/bbb-html5.nginx
#systemctl enable nginx

#systemctl enable tomcat9
#systemctl start tomcat9


#Ensure that IP is correct for configs of mediasoup, kurento, freeswitch
sudo sed -i "s/172.17.0.2/$LOCAL_IP/g" /usr/share/bigbluebutton/nginx/sip.nginx
systemctl restart nginx

FREESWITCH_CONFIG_PREV_CONTENT=$(sudo cat /opt/freeswitch/etc/freeswitch/vars.xml /opt/freeswitch/etc/freeswitch/sip_profiles/external.xml)
sudo sed -i "s/172.17.0.2/$LOCAL_IP/g" /opt/freeswitch/etc/freeswitch/vars.xml
sudo sed -i "s/172.17.0.2/$LOCAL_IP/g" /opt/freeswitch/etc/freeswitch/sip_profiles/external.xml

[ "$FREESWITCH_CONFIG_PREV_CONTENT" != "$(sudo cat /opt/freeswitch/etc/freeswitch/vars.xml /opt/freeswitch/etc/freeswitch/sip_profiles/external.xml)" ] && sudo systemctl restart freeswitch

WEBRTCSFU_CONFIG_PREV_CONTENT=$(cat /usr/local/bigbluebutton/bbb-webrtc-sfu/config/default.yml)
yq e -i '.mediasoup.webrtc.listenIps[0].ip = "0.0.0.0"' /usr/local/bigbluebutton/bbb-webrtc-sfu/config/default.yml
yq e -i ".mediasoup.webrtc.listenIps[0].announcedIp = \"$LOCAL_IP\"" /usr/local/bigbluebutton/bbb-webrtc-sfu/config/default.yml
yq e -i '.mediasoup.plainRtp.listenIp.ip = "0.0.0.0"' /usr/local/bigbluebutton/bbb-webrtc-sfu/config/default.yml
yq e -i ".mediasoup.plainRtp.listenIp.announcedIp = \"$LOCAL_IP\"" /usr/local/bigbluebutton/bbb-webrtc-sfu/config/default.yml
yq e -i '.kurento[0].ip = "$LOCAL_IP"' /usr/local/bigbluebutton/bbb-webrtc-sfu/config/default.yml
#yq e -i '.kurento[0].url = ws://127.0.0.1:8888/kurento' /usr/local/bigbluebutton/bbb-webrtc-sfu/config/default.yml
yq e -i ".freeswitch.ip = \"$LOCAL_IP\"" /usr/local/bigbluebutton/bbb-webrtc-sfu/config/default.yml
yq e -i ".freeswitch.sip_ip = \"$LOCAL_IP\"" /usr/local/bigbluebutton/bbb-webrtc-sfu/config/default.yml
[ "$WEBRTCSFU_CONFIG_PREV_CONTENT" != "$(cat /usr/local/bigbluebutton/bbb-webrtc-sfu/config/default.yml)" ] && sudo systemctl restart bbb-webrtc-sfu


#html5: set correct hostname
BBBHTML5_CONFIG_PREV_CONTENT=$(cat /etc/bigbluebutton/bbb-html5.yml)
sudo yq e -i ".public.pads.url = \"https://$(hostname -f)/pad\"" /etc/bigbluebutton/bbb-html5.yml
sudo yq e -i ".public.kurento.wsUrl = \"wss://$(hostname -f)/bbb-webrtc-sfu\"" /etc/bigbluebutton/bbb-html5.yml
[ "$BBBHTML5_CONFIG_PREV_CONTENT" != "$(cat /etc/bigbluebutton/bbb-html5.yml)" ] && sudo systemctl restart bbb-html5

#html5: make akka-apps read the same settings.yml of bbb-html5 project (using symlink)
BBBHTML5_USR_SHARE_CONFIG_PREV_CONTENT=$(cat /usr/share/meteor/bundle/programs/server/assets/app/config/settings.yml)
if [ -f /home/bigbluebutton/src/bigbluebutton-html5/private/config/settings.yml ]; then
  sudo ln -sfn /home/bigbluebutton/src/bigbluebutton-html5/private/config/settings.yml /usr/share/meteor/bundle/programs/server/assets/app/config/settings.yml
fi
[ "$BBBHTML5_USR_SHARE_CONFIG_PREV_CONTENT" != "$(cat /usr/share/meteor/bundle/programs/server/assets/app/config/settings.yml)" ] && sudo systemctl restart bbb-apps-akka


#bbb-web: set securitySalt

BBBWEB_CONFIG_PREV_CONTENT=$(cat /etc/bigbluebutton/bbb-web.properties)
sudo sed -i "/securitySalt=/d" /etc/bigbluebutton/bbb-web.properties || true

echo "securitySalt=$(sudo bbb-conf --salt | grep Secret: | cut -d ' ' -f 6)" | sudo tee -a /etc/bigbluebutton/bbb-web.properties
[ "$BBBWEB_CONFIG_PREV_CONTENT" != "$(cat /etc/bigbluebutton/bbb-web.properties)" ] && sudo systemctl restart bbb-web

#Avoid restarting BBB to speed up the environment startup
# /usr/bin/bbb-conf --restart

mkdir -p /usr/share/bbb-libreoffice
cd /usr/share/bbb-libreoffice
sudo docker build -t bbb-soffice docker/

echo "rc.local executed"
